# scripts/bambisleep-circuit-breaker.service
[Unit]
Description=BambiSleep.Chat Circuit Breaker Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=brandynette
Group=brandynette
WorkingDirectory=/home/brandynette/web/bambisleep.chat/js-bambisleep-chat
ExecStart=/usr/bin/node /home/brandynette/web/bambisleep.chat/js-bambisleep-chat/scripts/circuit-breaker-server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bambisleep-circuit-breaker

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/run /home/brandynette/web/bambisleep.chat/js-bambisleep-chat

# Environment
Environment=NODE_ENV=production
Environment=CIRCUIT_BREAKER_PORT=6970

[Install]
WantedBy=multi-user.target
