# scripts/bambisleep.service
[Unit]
Description=BambiSleep.Chat Node.js Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=brandynette
Group=brandynette
WorkingDirectory=/home/brandynette/web/bambisleep.chat/js-bambisleep-chat
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bambisleep-chat

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/run /home/brandynette/web/bambisleep.chat/js-bambisleep-chat

# Environment
Environment=NODE_ENV=production
Environment=SERVER_PORT=6969

[Install]
WantedBy=multi-user.target