# scripts/bambisleep.service
[Unit]
Description=BambiSleep.Chat Neural Network Control System
After=network.target

[Service]
Type=simple
User=brandynette
WorkingDirectory=/home/brandynette/web/bambisleep.chat/js-bambisleep-chat
ExecStart=/home/brandynette/.nvm/versions/node/v24.1.0/bin/node src/server.js
Restart=always
Environment=NODE_ENV=production
EnvironmentFile=-/home/brandynette/web/bambisleep.chat/js-bambisleep-chat/.env

[Install]
WantedBy=multi-user.target
